# scripts/sentinel-api.service
[Unit]
Description=Sentinel Intelligence CTI API
Documentation=https://github.com/sentinel-cti
After=network.target tailscaled.service
Wants=tailscaled.service

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/sentinel
EnvironmentFile=/home/ubuntu/sentinel/.env
ExecStart=/home/ubuntu/sentinel/venv/bin/uvicorn backend.api.main:app --host 0.0.0.0 --port 8000 --workers 1
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
